# examples/lxc-server/packet/etc/profile.d/ruby.sh
export PATH="$PATH:/var/lib/gems/1.8/bin"
export RUBYOPT="-Ku -rubygems"

